# templates/default/mysql.service.erb
# MySQL systemd service drop-in
# Maintained by Chef

[Service]
LimitNOFILE=<%= node["mysqlx"]["tunable"]["open_files_limit"] %>
